# conf/timetraveler-scan-config@.service
[Unit]
Description=Scan for updates to timetraveler configs across all users

[Service]
ExecStart=/usr/bin/timetraveler scan-config
Type=simple
