# pkg-src/DEBIAN/prerm
#!/bin/bash

if [ "$1" != "upgrade" ]; then
    find /etc/systemd/system/ -name 'timetraveler-*' | while read unit; do
        systemctl disable --now "$(basename "$unit")"
    done
fi

