# pkg-src/DEBIAN/prerm
#!/bin/bash

if [ "$1" != "upgrade" ]; then
    systemctl disable --now timetraveler-scan-config.timer
fi

